# templates/peslackbot.service.epp
<%- |
  String $slack_api_key,
| -%>
[Unit]
Description = peslackbot
After = pe-orchestration-services.service

[Service]
User = pe-puppet
Group = pe-puppet
Type = simple
Restart = always
Environment=SLACK_API_TOKEN=<%= $slack_api_key %>
WorkingDirectory=/opt/pe-slack-bot
ExecStart = /opt/puppetlabs/puppet/bin/puma
StandardOutput=syslog

[Install]
WantedBy = multi-user.target
